# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:pm2]
command=pm2 start /quran/start.js -i 4 --no-daemon
directory=/quran
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
auto_start=true
autorestart=true
